# cron/llm-nightly.timer.template
# Phase 3.8a-Mini systemd user timer template.
# See cron/llm-nightly.service.template for install instructions.
#
# OnCalendar fires 02:00 every day in the user's timezone.
# Persistent=true catches up if the timer was missed (WSL shutdown).
# RandomizedDelaySec spreads load if multiple matrices are ever scheduled.

[Unit]
Description=Schedule for autonomous-local-llm nightly experiment matrix
# NOTE: do NOT add `Requires=llm-nightly.service` here. systemd interprets
# `Requires` in a timer's [Unit] as "the service must be active whenever the
# timer is active", which causes `systemctl start <timer>` to pull-in the
# service immediately (verified 2026-06-01 fire test). The `Unit=` directive
# in [Timer] below is all that is needed to schedule the service.

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300
Unit=llm-nightly.service

[Install]
WantedBy=timers.target
